# big-mount/DEBIAN/postinst
#! /bin/sh

set -e

systemctl daemon-reload
systemctl enable bigmount.service 


exit 0
